# DEBIAN/prerm
#!/bin/sh
/etc/init.d/plexmediaserver stop
update-rc.d plexmediaserver remove
